# scripts/com.calf.mount-shares.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.calf.mount-shares</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/calf-mount-shares.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/calf-mount.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/calf-mount.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/admin</string>
        <key>USER</key>
        <string>admin</string>
    </dict>
</dict>
</plist>
